# supervisor/conf.d/xray.ini
[program:xray]
command=sh /xray/bin/xray.sh
directory=/xray
autostart=true
autorestart=true
redirect_stderr=true
priority=3
stdout_logfile=/tmp/xray.log