# supervisor/conf.d/xray.ini
[program:xray]
command=/xray/xray-core/xray -config /xray/xray-core/server.json
directory=/xray/xray-core
autostart=true
autorestart=true
redirect_stderr=true
priority=10
stdout_logfile=/tmp/xray.log